# templates/docker-oom-exporter.conf.j2
exec /usr/local/bin/docker_oom_exporter -output /var/lib/node_exporter/textfile_collector/docker_oom.prom

start on started docker
respawn
